# FILES.gluster/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:glusterd]
command=/usr/sbin/glusterd -N
